# misc/orbit.service
[Unit]
Description=orbit agent

[Service]
ExecStart=/usr/bin/orbit agent
Restart=always

[Install]
WantedBy=multi-user.target
